# systemd/armt-vpn-web.service
[Unit]
Description=ARMT VPN Web Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/armt-vpn
EnvironmentFile=/opt/armt-vpn/.env
ExecStart=/usr/bin/node /opt/armt-vpn/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
