# systemd/armt-vpn-bot.service
[Unit]
Description=ARMT VPN Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/armt-vpn
EnvironmentFile=/opt/armt-vpn/.env
ExecStart=/opt/armt-vpn/venv/bin/python /opt/armt-vpn/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
